# deployments/systemd/sdkwork-drive-admin-storage-api.service
[Unit]
Description=SDKWork Drive Admin Storage API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sdkwork
Group=sdkwork
EnvironmentFile=-/etc/sdkwork-drive/admin-storage-api.env
Environment=SDKWORK_DRIVE_ADMIN_STORAGE_API_BIND=0.0.0.0:18083
Environment=SDKWORK_DRIVE_RUNTIME_PROFILE=production
ExecStart=/usr/local/bin/sdkwork-router-storage-backend-api
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
