# deployments/systemd/sdkwork-drive-open-api.service
[Unit]
Description=SDKWork Drive Open API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sdkwork
Group=sdkwork
EnvironmentFile=-/etc/sdkwork-drive/open-api.env
Environment=SDKWORK_DRIVE_OPEN_API_BIND=0.0.0.0:18082
Environment=SDKWORK_DRIVE_RUNTIME_PROFILE=production
Environment=SDKWORK_DRIVE_DEPLOYMENT_PROFILE=cloud
ExecStart=/usr/local/bin/sdkwork-router-drive-open-api
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
