# deployments/systemd/sdkwork-drive-backend-api.service
[Unit]
Description=SDKWork Drive Backend API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sdkwork
Group=sdkwork
EnvironmentFile=-/etc/sdkwork-drive/backend-api.env
Environment=SDKWORK_DRIVE_BACKEND_API_BIND=0.0.0.0:18081
Environment=SDKWORK_DRIVE_RUNTIME_PROFILE=production
Environment=SDKWORK_DRIVE_DEPLOYMENT_PROFILE=cloud
ExecStart=/usr/local/bin/sdkwork-router-drive-backend-api
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
